# misc/cyp.service.template
[Unit]
Description=Control Your Player

[Service]
WorkingDirectory=$PWD
ExecStart=node .

[Install]
WantedBy=multi-user.target
