# misc/cyp.service.template
[Unit]
Description=Control Your Player

[Service]
WorkingDirectory=$PWD
ExecStart=/usr/bin/env node .

[Install]
WantedBy=default.target
